# lib/systemd/user/rsync-to-tank.timer
[Unit]
Description=Timer for rsync to tank

[Timer]
OnCalendar=daily

[Install]
WantedBy=default.target
